# crontab
0 3 * * * /app/syncvideo.py --output /output
0 23 * * * /app/app.py
